# Unit-Files/container-discovery.service
[Unit]
Description=Annonunce Myproject Database
After=etcd.service

[Service]
ExecStart=
ExecStop=
